# deb/bananian-etc_armhf/DEBIAN/postinst
#!/bin/bash
update-rc.d bananian-settings defaults
