# deb/bananian-etc_armhf/DEBIAN/postrm
#!/bin/bash
update-rc.d -f bananian-settings remove
